# src/doc/sifit.cron.d
0,30 * * * * www-data /usr/local/bin/php /var/www/html/cron/launch_processes.php > /tmp/output-sifit.txt 2>&1
